# init/myinit.conf
start on stopped rc RUNLEVEL=[2345]
stop on starting runlevel [016]

respawn
exec /root/myinit.sh
